# profile.d/200.npm-global-path.sh
# set PATH to include the npm-global folder owned by me
if [ -d "$HOME/.npm/bin" ] ; then
    export PATH="$HOME/.npm/bin:$PATH"
fi
