# dotfiles/profile.d/200.path-local-bin.sh
# initialize $PATH with existing bin directories

PATH="/bin"

if [ -d "/sbin" ] ; then 
   PATH="/sbin:$PATH"
fi

if [ -d "/usr/bin" ] ; then 
   PATH="/usr/bin:$PATH"
fi

if [ -d "/usr/sbin" ] ; then 
   PATH="/usr/sbin:$PATH"
fi

if [ -d "/usr/local/bin" ] ; then 
   PATH="/usr/local/bin:$PATH"
fi

if [ -d "/usr/local/sbin" ] ; then 
   PATH="/usr/local/sbin:$PATH"
fi

if [ -d "$HOME/bin" ] ; then
    PATH="$HOME/bin:$PATH"
fi

if [ -d "$HOME/.local/bin" ] ; then
    PATH="$HOME/.local/bin:$PATH"
fi
