# dotfiles/profile.d/205-path-snap.sh
# add the snap binary to the path
PATH="/snap/bin:$PATH"
